# s6-overlay/s6-rc.d/webapp/run
#!/bin/sh


# Inherit container environment from s6
ENV_DIR=""
[ -d /run/s6/container_environment ] && ENV_DIR="/run/s6/container_environment"
[ -d /var/run/s6/container_environment ] && ENV_DIR="/var/run/s6/container_environment"
if [ -n "$ENV_DIR" ]; then
    for f in "$ENV_DIR"/*; do
        [ -f "$f" ] && export "$(basename "$f")=$(cat "$f")"
    done
fi


RQLITE_URL="http://localhost:${RQLITE_HTTP_PORT:-4001}"
echo "Waiting for rqlite at ${RQLITE_URL}..."
until wget -q -O /dev/null "${RQLITE_URL}/readyz" 2>/dev/null; do
    sleep 1
done
echo "rqlite is ready."


cd /app
exec ./node_modules/.bin/next start -p ${PORT:-3000}