# s6-overlay/s6-rc.d/rqlite/run
#!/command/execlineb -P


backtick -n RQLITE_JOIN_FLAG {
    importas -D "" RQLITE_JOIN RQLITE_JOIN
    if { test -n "${RQLITE_JOIN}" }
    echo "-join=${RQLITE_JOIN}"
}


importas -D "node1" NODE_NAME NODE_NAME
importas -D "4001" RQLITE_HTTP_PORT RQLITE_HTTP_PORT
importas -D "4002" RQLITE_RAFT_PORT RQLITE_RAFT_PORT
importas -D "/data/rqlite" RQLITE_DATA_DIR RQLITE_DATA_DIR


exec /usr/local/bin/rqlited
    -node-id=${NODE_NAME}
    -http-addr=0.0.0.0:${RQLITE_HTTP_PORT}
    -raft-addr=0.0.0.0:${RQLITE_RAFT_PORT}
    ${RQLITE_JOIN_FLAG}
    ${RQLITE_DATA_DIR}