# s6-overlay/s6-rc.d/ws-server/run
#!/usr/bin/with-contenv sh
set -eu

RQLITE_HTTP_PORT="${RQLITE_HTTP_PORT:-4001}"
RQLITE_URL="http://127.0.0.1:${RQLITE_HTTP_PORT}"

echo "[ws-server] Waiting for rqlite at ${RQLITE_URL}/readyz ..."
until wget -q -O /dev/null "${RQLITE_URL}/readyz" 2>/dev/null; do
    sleep 1
done
echo "[ws-server] rqlite is ready. Starting ws-server..."

cd /app
exec node --input-type=module /app/ws-server.mjs